# contrib/systemd/xembsni.service
[Unit]
Description=XEmbed to StatusNotifierItem tray bridge
Documentation=https://github.com/jmylchreest/xembsni
# Bind to the graphical session so it starts/stops with the compositor
# (niri, Hyprland, etc. that integrate with systemd's graphical-session.target).
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/xembsni
Restart=on-failure
RestartSec=2
# Logs go to the journal via stderr; override verbosity with e.g.
#   systemctl --user edit xembsni.service
# and add:  [Service]\n Environment=RUST_LOG=debug
Environment=RUST_LOG=info

[Install]
WantedBy=graphical-session.target
